# twistedtv-video-server.service
[Unit]
Description=TwistedTV Video Streaming Server
After=network.target

[Service]
Type=exec
User=twistedtv
WorkingDirectory=/home/twistedtv/cinema-chat/twistedtv-video-server
Environment="PATH=/home/twistedtv/cinema-chat/twistedtv-server/venv/bin:/usr/bin:/usr/local/bin"
ExecStartPre=/bin/bash -c 'lsof -ti:9000 | xargs -r kill -9 || true'
ExecStart=/home/twistedtv/cinema-chat/twistedtv-server/venv/bin/python streaming_server.py
Restart=on-failure
RestartSec=10
StandardOutput=append:/tmp/twistedtv-video-server.log
StandardError=append:/tmp/twistedtv-video-server.log

[Install]
WantedBy=multi-user.target
